# deploy/mattrmindr.service
# MattrMindr Case Management System — systemd service
#
# Before installing, update the paths below to match your install location:
#   - WorkingDirectory
#   - EnvironmentFile
#   - StandardOutput
#   - StandardError
#
# The server loads .env automatically via dotenv, so EnvironmentFile
# is a safety net — either one is sufficient.

[Unit]
Description=MattrMindr Case Management System
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/path/to/MattrMindr-Demo
EnvironmentFile=/path/to/MattrMindr-Demo/.env
ExecStart=/usr/bin/node server/index.js
Restart=on-failure
RestartSec=5
StandardOutput=append:/path/to/MattrMindr-Demo/logs/out.log
StandardError=append:/path/to/MattrMindr-Demo/logs/error.log

[Install]
WantedBy=multi-user.target
